# podvm-mkosi/mkosi.skeleton/etc/profile.d/20-ssh-banner.sh
if [[ -n ${SSH_CONNECTION} && $- == *i* ]] ; then
    neofetch
fi
